# lab05/solution/etc/init.d/S99server
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
    echo "Startuję serwer Python Tornado w tle..."
    cd /usr/server
    # Uruchomienie w tle; przekierowanie logów do /var/log/server.log
    nohup python3 server.py 5001 192.168.144.15 >> /var/log/server.log 2>&1 &
    # Zapisz PID do pliku, żeby później móc zabić właściwy proces
    echo $! > /var/run/server.pid
    echo "Serwer uruchomiony, PID=$(cat /var/run/server.pid)."
}

stop() {
    echo "Zatrzymuję serwer Python Tornado..."
    if [ -f /var/run/server.pid ]; then
        PID=$(cat /var/run/server.pid)
        kill "$PID" 2>/dev/null
        rm -f /var/run/server.pid
        echo "Serwer (PID=$PID) zatrzymany."
    else
        echo "Plik /var/run/server.pid nie istnieje – serwer prawdopodobnie nie był uruchomiony."
    fi
}

